# ami/mkosi.extra/usr/lib/systemd/system/rke2-bootstrap.path
[Unit]
Description=Watch for RKE2 bootstrap env file

[Path]
PathExists=/run/rke2-bootstrap/env

[Install]
WantedBy=multi-user.target
